# deploy/systemd/classflow-cleanup.service
[Unit]
Description=ClassFlow Temporary Directory Cleanup

[Service]
Type=oneshot
WorkingDirectory=/opt/classflow
EnvironmentFile=/etc/classflow/backend.env
ExecStart=/opt/classflow/scripts/classflow-cleanup-temp.sh
